# linux_memory_and_processes/homework_service.service
[Unit]
Description=Simple homework service writing date

[Service]
Type=exec
ExecStart=/usr/local/bin/homework_service.sh
Restart=on-failure